# packaging/prerm
#!/bin/sh
set -e
if [ "$1" = remove ] && command -v systemctl >/dev/null 2>&1; then
  systemctl disable minicontainer-reconcile.service >/dev/null 2>&1 || true
fi
